# debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
    systemctl stop debounced.service || true
    systemctl disable debounced.service || true
fi
exit 0